# config/systemd/user/keybase.service.d/no-debug.conf
[Service]
ExecStart=
ExecStart=/usr/bin/keybase --no-debug service
